# systemd/camera.service
[Unit]
Description=Camera

[Service]
WorkingDirectory=/opt/picam
ExecStart=/usr/bin/python3 -m mncam

[Install]
WantedBy=multi-user.target
